# overlay/usb-role-manager.service
[Unit]
Description=Auto USB role manager (host/charger)
After=local-fs.target systemd-udevd.service
Wants=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/usb-mode-switch.sh auto
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
