# overlay/usb-otg-host.service
[Unit]
Description=Auto USB mode switch (host/charger detection)
After=sys-devices-platform-ff740000.usb2\x2dphy.device

[Service]
Type=simple
ExecStart=/usr/local/bin/usb-mode-switch.sh auto
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
